# deploy/systemd/koordinator-celery.service
[Unit]
Description=ГИС КООРДИНАТОР Celery Worker
After=network.target redis.service mysql.service
Wants=redis.service

[Service]
Type=simple
User=tonojkeee
Group=tonojkeee
WorkingDirectory=/home/tonojkeee/projects/main/backend
Environment="PATH=/home/tonojkeee/projects/main/backend/venv/bin"
EnvironmentFile=/home/tonojkeee/projects/main/backend/.env

ExecStart=/home/tonojkeee/projects/main/backend/venv/bin/celery -A app.core.celery_app worker --loglevel=info --concurrency=4

# Graceful shutdown
TimeoutStopSec=300
KillMode=mixed

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=koordinator-celery

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
