# deploy/systemd/koordinator-smtp.service
[Unit]
Description=ГИС КООРДИНАТОР SMTP Server
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=tonojkeee
Group=tonojkeee
WorkingDirectory=/home/tonojkeee/projects/main/backend
Environment="PATH=/home/tonojkeee/projects/main/backend/venv/bin"
EnvironmentFile=/home/tonojkeee/projects/main/backend/.env

# Run SMTP server standalone (optional, included in main app by default)
ExecStart=/home/tonojkeee/projects/main/backend/venv/bin/python -m app.modules.email.smtp_standalone

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=koordinator-smtp

# Security
NoNewPrivileges=true
PrivateTmp=true
# Allow binding to port 2525
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
